# systemd/vigia.service
[Unit]
Description=RondaVizinho (daemon)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m vigia
WorkingDirectory=/opt/vigia
Environment=VIGIA_CONFIG=/var/lib/vigia/config.json
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
